# debian/postinst
#!/bin/bash
set -e

# Post-installation script for onefinity-firmware

case "$1" in
    configure)
        # Create system user if it doesn't exist
        if ! id -u onefinity >/dev/null 2>&1; then
            useradd -r -s /usr/sbin/nologin -d /var/lib/onefinity -c "OneFinity CNC Service" onefinity
        fi
        
        # Install Python packages not available in Debian repositories
        # sockjs-tornado is no longer available as a system package in Debian 12+
        echo "Installing Python dependencies via pip..."
        pip3 install --break-system-packages sockjs-tornado || {
            echo "Warning: Failed to install sockjs-tornado via pip"
            exit 1
        }
        
        # Reload systemd daemon
        systemctl daemon-reload
        
        # Enable but don't start the service (user should configure first)
        systemctl enable onefinity.service || true
        
        # Set permissions for application directory
        chown -R root:root /opt/onefinity
        chmod -R 755 /opt/onefinity/bin
        
        # Create log directory with proper permissions
        mkdir -p /var/log/onefinity
        chown onefinity:onefinity /var/log/onefinity
        chmod 755 /var/log/onefinity
        
        # Create state directory with proper permissions
        mkdir -p /var/lib/onefinity
        chown onefinity:onefinity /var/lib/onefinity
        chmod 755 /var/lib/onefinity
        
        echo "OneFinity firmware installed successfully."
        echo "To start the service, run: sudo systemctl start onefinity"
        ;;
    
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
